# zoterobck.timer
[Unit]
Description=Backup Zotero using git

[Timer]
OnCalendar=daily
Persistent=true
Unit=zoterobck.service

[Install]
WantedBy=timers.target

