# zoterobck.service
[Unit]
Description="A script to backup all home dirs to usb drive"

[Service]
ExecStart=/home/jcuenod/Zotero/backup.sh
