# zoterobck.service
[Unit]
Description="Execute Zotero backup script which uses git to version control research. Booya!"

[Service]
ExecStart=%h/Zotero/backup.sh
