# zoterobck.service
[Unit]
Description="A script to backup all home dirs to usb drive"

[Service]
ExecStart=~/Zotero/backup.sh
